# packaging/t1bar-scenes.service
[Unit]
Description=t1bar-studio - drives the T1 Touch Bar from your SCENES config (Scenes/Layer-Loom engine)
Documentation=https://github.com/AJ-dev-i60/t1-touchbar
# Up before the greeter so the bar + config 2 are ready at login (no cycle: do NOT
# also add After=multi-user.target while WantedBy=multi-user.target).
Before=display-manager.service

[Service]
Type=simple
# Runs as root (USB control, briefly unloads the apple touchbar modules, uinput).
# Reads your SCENES config and HOT-RELOADS it, so edits from Scene Home appear on the
# bar live (~1s). Resolves the active scene from live machine state; motion ticks at
# <=30fps, idle is event-driven. Media/volume bridge to the logged-in user's session.
ExecStart=/usr/bin/env t1bar scene-run -c CONFIG_PATH
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
