# packaging/t1touchbar-strip.service
[Unit]
Description=T1 Touch Bar control strip (t1-touchbar)
Documentation=https://github.com/AJ-dev-i60/t1-touchbar
# Be up before the display manager so the bar + config 2 are ready for the greeter.
# (Do NOT use After=multi-user.target here: combined with WantedBy below and the
# camera unit's ordering it forms a cycle and systemd drops a service at boot.)
Before=display-manager.service

[Service]
Type=simple
# Runs as root: needs USB control + briefly unloading the kernel touchbar modules.
# Media/volume are bridged to the logged-in user's session automatically.
ExecStart=/usr/bin/env t1touchbar-strip --welcome-once
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
